# asuman-memory.service.example
# ╔═══════════════════════════════════════════════════════════════╗
# ║  CHANGE all /opt/your-org/ paths to your install location   ║
# ║  CHANGE EnvironmentFile to your .env path                   ║
# ║  CHANGE ReadWritePaths to your DB directory                  ║
# ╚═══════════════════════════════════════════════════════════════╝

[Unit]
Description=Agent Memory API
After=network.target

[Service]
Type=simple
User=root

# >>> CHANGE these paths <<<
WorkingDirectory=/opt/your-org/asuman-memory
ExecStart=/opt/your-org/asuman-memory/.venv/bin/python -m agent_memory

Restart=always
RestartSec=5

# >>> CHANGE to your env file path <<<
# This file should contain OPENROUTER_API_KEY, AGENT_MEMORY_API_KEY, etc.
EnvironmentFile=/opt/your-org/asuman-memory/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=asuman-memory

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ProtectHome=read-only
# >>> CHANGE to your DB directory <<<
ReadWritePaths=/home/your-user/.asuman
CapabilityBoundingSet=
UMask=0077

[Install]
WantedBy=multi-user.target
